# deploy/ticker-tide-bot.service
[Unit]
Description=Ticker Tide Telegram Bot
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=PROJECT_DIR
ExecStart=PROJECT_DIR/.venv/bin/python scripts/run_bot.py
EnvironmentFile=PROJECT_DIR/.env
Restart=always
RestartSec=5
StandardOutput=append:PROJECT_DIR/logs/bot.log
StandardError=append:PROJECT_DIR/logs/bot.log

[Install]
WantedBy=multi-user.target
